# scripts/http-monitor.service
[Unit]
Description=HTTP Traffic Monitor (tcpdump + python parser)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'exec tcpdump -i eth0 -A -s 0 -l "port 3000" 2>/dev/null | exec python3 /opt/http-monitor/http_monitor.py /var/log/http_monitor.json'
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/http-monitor

# 安全限制
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log
PrivateTmp=true

# 资源限制
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
